# scripts/pi/templates/parsebox-kiosk.service.template
[Unit]
Description=ParseBox kiosk HTTP server
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m http.server 4174 --directory __REPO_DIR__/kiosk
Restart=on-failure

[Install]
WantedBy=default.target
